# sietchtabr/rc.local
for ruleset in vpns lans lan-based-services; do
  /sbin/pfctl -a ${ruleset} -f /etc/pf.conf.${ruleset}
done
