# sietchtabr/rc.local
for ruleset in vpns lans; do
  /sbin/pfctl -a ${ruleset} -f /etc/pf.conf.${ruleset}
done
